# mmm.conf
[program:mmm]
command=python3 /home/curator2/mmm/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator2/mmm/test.err.log
stdout_logfile=/home/curator2/mmm/test.out.log